# root/etc/s6-overlay/s6-rc.d/parrot/run
#!/command/execlineb -P
npm start --prefix /parrot